# debian/di-live.service
[Unit]
Description=di-live: TurnKey Linux installer
After=getty@tty8.service
Before=rescue.service
ConditionKernelCommandLine=di-live

[Service]
Type=idle
StandardInput=tty-force
TTYPath=/dev/tty8
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
ExecStart=/etc/init.d/di-live start
ExecStop=/etc/init.d/di-live stop

[Install]
WantedBy=rescue.target
